# services/test2/run
#!/bin/bash
bundle exec rawit notifier start
trap 'exit 0' INT TERM
while true
do
sleep 1
done
