# services/test2/run
#!/bin/bash
bundle exec rawit notifier start
while true
do
sleep 1
done
